# deploy/backup.timer
[Unit]
Description=MemBind Database Backup Timer (hourly)

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
